# systemd/sipquest-box.service
[Unit]
Description=SipQuest Box Runtime
After=sound.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/sipquest
EnvironmentFile=/opt/sipquest/.env
ExecStart=/opt/sipquest/.venv/bin/sipquest-voice-box
Restart=always
RestartSec=3
RuntimeDirectory=sipquest
StateDirectory=sipquest
SupplementaryGroups=audio gpio i2c spi input

[Install]
WantedBy=multi-user.target
